# server/docrec.service
[Unit]
Description=DocRec Flask app (docrec)
After=network.target

[Service]
User=USER_PLACEHOLDER
WorkingDirectory=/home/USER_PLACEHOLDER/docrec/server
Environment=PORT=8080
ExecStart=/home/USER_PLACEHOLDER/docrec/.venv/bin/gunicorn -b 0.0.0.0:8080 app:app --workers 2
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target